# etc/grimoire.service
[Unit]
Description=Grimoire multi-GPU llama.cpp gateway
Wants=network-online.target docker.service
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=simple
TimeoutStartSec=900
TimeoutStopSec=120
Restart=on-failure
RestartSec=10
User=yeowool
WorkingDirectory=/home/yeowool/grimoire

ExecStartPre=-/usr/bin/docker compose down

ExecStart=/usr/bin/docker compose up

ExecStop=/usr/bin/docker compose down

[Install]
WantedBy=multi-user.target
